# supervisor/alarm.conf
[program:alarm]
command=python3 test.py
directory=/home/hebron/rpi-alarm
autostart=true
autorestart=unexpected
user=hebron
stderr_logfile=syslog
stdout_logfile=syslog
